# services/data/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ─── JBT DATA API (FastAPI/uvicorn) ───────────────
[program:jbt-data-api]
command=python -m uvicorn src.main:app --host 0.0.0.0 --port 8105
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=10
stdout_logfile=/data/logs/api_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/api_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
environment=PYTHONPATH="/app",DATA_STORAGE_ROOT="/data",JBT_DATA_LOG_DIR="/data/logs"

; ─── JBT DATA Scheduler (APScheduler 24h) ────────
[program:jbt-data-scheduler]
command=python -m src.scheduler.data_scheduler
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=30
stdout_logfile=/data/logs/scheduler_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/scheduler_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
environment=PYTHONPATH="/app",DATA_STORAGE_ROOT="/data",JBT_DATA_LOG_DIR="/data/logs"

; ─── JBT DATA Health Monitor (periodic, not a daemon) ───
[program:jbt-data-health]
command=python -m src.health.health_check
directory=/app
autostart=true
autorestart=unexpected
startsecs=0
exitcodes=0
startretries=3
stopwaitsecs=10
stdout_logfile=/data/logs/health_stdout.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
stderr_logfile=/data/logs/health_stderr.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=2
environment=PYTHONPATH="/app",DATA_STORAGE_ROOT="/data",JBT_DATA_LOG_DIR="/data/logs"
